# docker/s6-rc.d/api-server/finish
#!/command/with-contenv sh
echo "API Server stopped"
