# docker/s6-rc.d/api-server/run
#!/command/with-contenv sh
exec 2>&1
echo "[api-server] Starting api-server serve..."
exec api-server serve
