# docker/s6-rc.d/api-server/run
#!/command/with-contenv sh
exec api-server serve
